# roles/redshift/templates/redshift.service
[Unit]
Description=Redshift

[Service]
ExecStart=/usr/local/bin/redshift
Restart=on-failure

[Install]
WantedBy=default.target
