# roles/redshift/templates/redshift.service
[Unit]
Description=Redshift

[Service]
ExecStart=/usr/bin/redshift -l 48.8:-122.5
Restart=on-failure

[Install]
WantedBy=default.target
